# systemd/wgaw/service/system-stats.timer
[Unit]
Description=Custom Timer

[Timer]
# OnBootSec=1min
OnBootSec=10
# OnUnitActiveSec=1h
OnUnitActiveSec=10

[Install]
WantedBy=timers.target